# base/.config/systemd/user/rmpc-recent-playlist.service
[Unit]
Description=Update MPD Recently Added Playlist
After=mpd.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rmpc_recent_playlist.sh
