# base/.config/systemd/user/rmpc-recent-playlist.timer
[Unit]
Description=Hourly timer for MPD Recently Added Playlist update

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
